# deploy/systemd/vastoria-website.service
[Unit]
Description=Vastoria Website
After=network-online.target vastoria-cloud-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/vastoria/apps/website
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
Environment=VASTORIA_ENV=production
Environment=CLOUD_API_URL=http://127.0.0.1:18430
Environment=NEXT_PUBLIC_CLOUD_API_URL=https://vastoria.online/api
Environment=DOWNLOADS_BASE_URL=https://vastoria.online/downloads
Environment=NEXT_PUBLIC_DOWNLOADS_URL=https://vastoria.online/downloads
ExecStart=/usr/bin/node .next/standalone/apps/website/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
